# rootfs/etc/local.d/local.start
#!/bin/sh

HOSTID=$(dmidecode | grep -w UUID | sed "s/^.UUID\: //g")

mount /dev/disk/by-label/CFG /etc/cfg 2>/dev/null
if [ -x /etc/cfg/${HOSTID} -a ! -z "${HOSTID}" ]; then
	/etc/cfg/${HOSTID}
	umount /etc/cfg
elif [ -x "/etc/rc.local.d/${HOSTID}" -a ! -z "${HOSTID}" ]; then
	/etc/rc.local.d/${HOSTID}
#elif try config download
else
	/etc/rc.local.d/rc.default
	/etc/rc.local.d/rc.fwstop
	/etc/rc.local.d/rc.fwstart
fi

sync || true
